# hopingswings4415.conf
[program:hopingswings4415]
command=python3 /home/curator00/hopingswings4415/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator00/hopingswings4415/test.err.log
stdout_logfile=/home/curator00/hopingswings4415/test.out.log